# collectors/systemd/calendar-collector.service
[Unit]
Description=Calendar Data Collector
After=network.target

[Service]
Type=oneshot
User=monitoring
Group=monitoring
WorkingDirectory=/home/olimex/sequoia_fabrica_blog
ExecStart=/usr/bin/node collectors/calendar-collector.js
StandardOutput=journal
StandardError=journal