# collectors/systemd/esp-logger.service
[Unit]
Description=ESP32 Serial Logger
After=network.target

[Service]
Type=simple
User=monitoring
Group=dialout
WorkingDirectory=/opt/sequoia_fabrica_blog
ExecStart=/usr/bin/python3 esp_logger.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
