# collectors/systemd/esp-logger.service
[Unit]
Description=ESP32 Battery Monitor Serial Logger
After=network.target

[Service]
Type=simple
User=monitoring
Group=dialout
ExecStart=/usr/bin/python3 /opt/sequoia_fabrica_blog/collectors/esp_logger.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ensure log directory exists
RuntimeDirectory=esp_logger
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
